# dot/systemd/tq.timer
[Unit]
Description=Task Queue Timer

[Timer]
OnBootSec=5min
OnActiveSec=0
OnUnitInactiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
